# linux/etc/systemd/user/rclone_onedriver.service
[Unit]
Description=One Drive (rclone)
AssertPathIsDirectory=/home/cccn/mnt/onedriver

[Service]
Type=simple
ExecStart=/bin/rclone mount onedriver: /home/cccn/mnt/onedriver
ExecStop=/bin/systemd-umount -u /home/cccn/mnt/onedriver
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
